# etc/service/dockergen/run
#!/bin/sh -e

exec /usr/local/bin/docker-gen -config /etc/docker-gen.cfg 2>&1
